# .config/systemd/user/ff-backup.service
[Unit]
Description=Firefox backup

[Service]
ExecStart=/home/nimrod/Documents/bin/ff-backup
