# .config/systemd/user/ff-backup.timer
[Unit]
Description=Firefox backup

[Timer]
OnCalendar=09:00
OnCalendar=21:00

[Install]
WantedBy=default.target
