# .config/autostart/dropbox.py.desktop
[Desktop Entry]
Exec=/home/nimrod/dropbox.py start
Icon=system-run
Path=
Terminal=false
Type=Application
